# roles/gunicorn/templates/gunicorn.socket.j2
[Unit]
Description={{ socket_file_name }} socket

[Socket]
ListenStream=/run/{{ socket_file_name }}.sock

[Install]
WantedBy=sockets.target